# packaging/server/rch-server.service
[Unit]
Description=Reticulum Community Hub Rust Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=rch
Group=rch
EnvironmentFile=/etc/rch/rch-server.env
ExecStart=/opt/rch/bin/r3akt-rch-server --bind ${RCH_BIND} --db-path ${RCH_DB_PATH} --ui-dist-path ${RCH_UI_DIST_PATH}
Restart=on-failure
RestartSec=5
StateDirectory=rch
LogsDirectory=rch

[Install]
WantedBy=multi-user.target
